# supervisord.conf
[program:example-app]
command=/root/example-app/example-app
